# dev-environment/config/abeone-orchestrator.service
[Unit]
Description=THE ONE SYSTEM Service Orchestrator
After=network.target docker.service
Wants=network.target docker.service

[Service]
Type=simple
User=user
Group=user
WorkingDirectory=/home/user/AbeOne_Master/dev-environment
ExecStart=/usr/bin/python3 /home/user/AbeOne_Master/dev-environment/scripts/orchestrator.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

